# printer-worker/cordial-printer.service
[Unit]
Description=Cordial KDS Printer Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/cordial-kds/printer-worker
EnvironmentFile=/home/pi/cordial-kds/printer-worker/.env
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
